# zebrakow.slice
[Unit]
Description=Default slice unit for zebrakow Job Scheduler
Documentation=https://github.com/bl33pbl0p/zebrakow/wiki
Before=slices.target

[Slice]
CPUAccounting=true
MemoryAccounting=true
TasksAccounting=true
IOAccounting=true
IPAccounting=true

